# debian/phoenixdkim.service
[Unit]
Description=PhoenixDKIM Milter
Documentation=man:phoenixdkim(8) man:phoenixdkim.conf(5) man:phoenixdkim-lua(3) man:phoenixdkim-genkey(8) man:phoenixdkim-genzone(8) man:phoenixdkim-testkey(8) man:phoenixdkim-testmsg(8)
After=network-online.target nss-lookup.target syslog.target
Wants=network-online.target

# Best-effort ordering so phoenixdkim is ready before a local MTA starts
# accepting mail.  The authoritative way to guarantee ordering is on the
# consumer side, e.g. "After=phoenixdkim.service" (and "Wants=phoenixdkim.service")
# in the MTA unit.
Before=postfix.service sendmail.service exim.service exim4.service opensmtpd.service

[Service]
Type=notify
RuntimeDirectory=phoenixdkim
RuntimeDirectoryMode=0750
ExecStart=/usr/sbin/phoenixdkim -f -x /etc/phoenixdkim/phoenixdkim.conf
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure

# phoenixdkim feeds this watchdog from a dedicated thread at half the interval
# below (i.e. every 15s).  systemd restarts the service if no keep-alive
# arrives within WatchdogSec, catching a wholly hung process.
WatchdogSec=30
User=phoenixdkim
Group=phoenixdkim
UMask=007

# As the file system is mounted read-only with ProtectSystem=strict, adding
# /var/spool/postfix/phoenixdkim by default makes it possible to use phoenixdkim
# with postfix without further configuration.
ReadWritePaths=-/var/spool/postfix/phoenixdkim

# node_exporter textfile collector directory (Debian: prometheus-node-exporter).
# If MetricsFile points elsewhere, add the correct path via:
#   systemctl edit phoenixdkim.service
ReadWritePaths=-/var/lib/prometheus/node-exporter

# Hardening - see systemd.exec(5)
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=true
PrivateMounts=true
PrivateTmp=true
PrivateUsers=true
ProcSubset=pid
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
RemoveIPC=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~ @privileged @resources

# These lines prevent phoenixdkim from executing any binary other than itself,
# which is a strong mitigation against an attacker spawning a shell.
# IMPORTANT: if you use ReportAddress or ReportCommand in phoenixdkim.conf,
# phoenixdkim will popen() a mail submission binary (sendmail by default).
# In that case you must add the binary's path, e.g.:
#   ExecPaths=/usr/sbin/phoenixdkim /usr/lib /usr/sbin/sendmail
# or override via: systemctl edit phoenixdkim.service
NoExecPaths=/
ExecPaths=/usr/sbin/phoenixdkim /usr/lib

[Install]
WantedBy=multi-user.target
